# phosphor-ted-dbus-interface.service.in
[Unit]
Description=Ted sensors

[Service]
Restart=always
ExecStart=@bindir@/phosphor-ted-dbus-interface
Type=dbus
BusName=xyz.openbmc_project.TedDbusInterface
SyslogIdentifier=phosphor-ted-dbus-interface

[Install]
WantedBy=multi-user.target
